# crontab
*/30 * * * * root /srv/Festa-crawling/app/festcp/spider.py >> /srv/Festa-crawling/cron.log 2>&1
